# etc/cron.d/cherrypeak

TZ=America/Denver

## cleanup
 0 20 * * * pi /usr/bin/find /home/pi/tmp -name "*.jpg" -mtime +1 -delete
30 20 * * * pi /usr/bin/find /home/pi/arc -name "*.jpg" -mtime +7 -delete

## monitoring
3,13,23,33,43,53 * * * * pi /home/pi/bin/monitor.pl upload

## network check
9,19,29,39,49,59 * * * * root /sbin/ifup wlan1

