# etc/cron.d/ewcd_power

TZ=America/Denver

## self-healing
58 * * * * root /sbin/ifup wlan0

