# etc/cron.d/network

TZ=America/Denver

## network check
0,19,29,39,49,59 * * * * root /sbin/ifup wlan0
#0,19,29,39,49,59 * * * * root /sbin/ifup wlan1

