# etc/cron.d/openvpn

TZ=America/Denver

0 * * * * root /usr/sbin/service openvpn restart

