# etc/cron.d/shoot

TZ=America/Denver

* * * * * pi /home/pi/bin/shoot.pl >> /home/pi/log/shoot

