# etc/cron.d/shoot

TZ=America/Denver

* * * * * pi /home/pi/bin/new.pl >> /home/pi/log/new

