# etc/cron.d/status

TZ=America/Denver

*/10 * * * * pi /home/pi/bin/status.pl

## cleanup
9 14 * * * pi /usr/bin/find /home/pi/status -name "*.txt" -mtime +7 -delete

