# etc/cron.d/timelapse

TZ=America/Denver

1 * * * * pi /home/pi/bin/timelapse_hourly.pl

