# etc/cron.d/upload_05m

TZ=America/Denver

*   * * * * pi /home/pi/bin/shoot.pl
*/5 * * * * pi /home/pi/bin/upload.pl

