# etc/cron.d/upload_02m

TZ=America/Denver

*   * * * * pi /home/pi/bin/shoot.pl
*/2 * * * * pi /home/pi/bin/upload.pl

